# muninn.service
[Unit]
Description=Muninn Telegram Bot
After=network.target

[Service]
User=llumetti
WorkingDirectory=/home/llumetti/.local/projects/Muninn
ExecStart=/usr/bin/python3 /home/llumetti/.local/projects/Muninn/muninn_bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 